# boot-stuff/timeshift/timeshift-monthly.timer
[Unit]
Description=Timeshift monthly snapshot timer

[Timer]
OnCalendar=*-*-01 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
