# boot-stuff/timeshift/timeshift-monthly.service
[Unit]
Description=Timeshift monthly snapshot

[Service]
Type=oneshot
ExecStart=/usr/bin/timeshift --create --comments "monthly snapshot" --tags M
