# boot-stuff/timeshift/timeshift-weekly.service
[Unit]
Description=Timeshift weekly snapshot

[Service]
Type=oneshot
ExecStart=/usr/bin/timeshift --create --comments "weekly snapshot" --tags W
